# app-admin/clockwork/files/clockd.initd
#!/sbin/openrc-run
# Distributed under the terms of the GNU General Public License v2

command=/usr/sbin/clockd
CLOCKD_CONF=${CLOCKD_CONF:-/etc/clockwork/clockd.conf}
commad_args="-c \"${CLOCKD_CONF}\""
pidfile=${CLOCKD_PIDFILE:-/run/clockd.pid}

depend() {
	need net
	use  dns logger
}
