# app-admin/clockwork/files/meshd.initd
#!/sbin/openrc-run
# Distributed under the terms of the GNU General Public License v2

command=/usr/sbin/meshd
MESHD_CONF=${MESHD_CONF:-/etc/clockwork/meshd.conf}
commad_args="-c \"${MESHD_CONF}\""
pidfile=${MESHD_PIDFILE:-/run/meshd.pid}

depend() {
	need net
	use  dns logger
}
