# www-apps/gitorious/files/crontab
* * * * * cd /var/www/gitorious/site && bundle exec /usr/bin/rake ultrasphinx:index RAILS_ENV="production"
