# www-apps/gitorious/files/gitorious-poller
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/dev-ruby/stompserver/files/stompserver,v 1.1 2009/03/25 17:02:32 caleb Exp $

GIT_DAEMON="/var/www/gitorious/site/script/poller"
GIT_USER="git"

start() {
        ebegin "Starting Gitorious poller"
        su - ${GIT_USER} -c "${GIT_DAEMON} start"
        eend $? "Failed to start Gitorious git-daemon"
}

stop() {
        ebegin "Shutting down Gitorious poller"
        su - ${GIT_USER} -c "${GIT_DAEMON} stop"
        eend $?
}

