# www-apps/gitorious/files/gitorious-git-daemon
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/dev-ruby/stompserver/files/stompserver,v 1.1 2009/03/25 17:02:32 caleb Exp $

PIDFILE=/var/run/git-daemon.pid
GIT_DAEMON="/var/www/gitorious/site/script/run-git-daemon"
GIT_USER="git"
GIT_GROUP="git"
PIDFILE="/var/www/gitorious/site/log/git-daemon.pid"

depend() {
        need net
        use logger
}

start() {
        ebegin "Starting Gitorious git-daemon"
	${GIT_DAEMON} start ${GIT_USER} ${GIT_GROUP}
        eend $? "Failed to start Gitorious git-daemon"
}

stop() {
        ebegin "Shutting down Gitorious git-daemon"
	${GIT_DAEMON} stop ${GIT_USER} ${GIT_GROUP}
        eend $?
}

